# nqptp.service.in
[Unit]
Description=NQPTP -- Not Quite PTP
Wants=network-online.target
After=network.target network-online.target
Before=shairport-sync.service

[Service]
ExecStart=@prefix@/bin/nqptp
User=nqptp
Group=nqptp

[Install]
WantedBy=multi-user.target
